# etc/supervisord.d/scratchsetup.ini
[program:sratchsetup]
command=/opt/qnib/bin/scratchsetup.sh
autorestart=false
redirect_stderr=true
stdout_logfile=syslog
