# etc/supervisord.d/slurm_update.ini
[program:slurm_update]
command=consul-template -consul consul.service.consul:8500 -template "/etc/consul-templates/slurm/slurm.conf.tmpl:/usr/local/etc/slurm.conf:/opt/qnib/slurm/bin/reread.sh"
redirect_stderr=true
stdout_logfile=syslog
