# etc/supervisord.d/usersetup.ini
[program:usersetup]
command=/opt/qnib/bin/usersetup.sh
autorestart=false
redirect_stderr=true
stdout_logfile=syslog
